# agent/packaging/ibmc-kvm-agent.service
[Unit]
Description=iBMC KVM Linux Agent
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/ibmc-kvm-agent serve
EnvironmentFile=-%h/.config/ibmc-kvm-agent/environment
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
UMask=0077

[Install]
WantedBy=graphical-session.target
